# services/MacArthurN2K/run
#!/bin/sh
# MacArthur HAT – NMEA 2000 to dbus bridge (generated by setup)
exec /usr/bin/python3 "SCRIPT_PATH_PH"
